# apm/supervisor/elasticsearch.ini
[program:elasticsearch]
command=/data/elk/elasticsearch-7.5.1/bin/elasticsearch
directory=/data/elk/elasticsearch-7.5.1
user=elk
numprocs=1
priority=1
autostart=true
startsecs=10
satrtretries=3
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
stdout_logfile=/var/log/supervisor/supervisor_elasticsearch.log
