# apm/supervisor/esalert.ini
[program:es_alert]
command=python -m elastalert.elastalert --verbose --rule example_rules/apm_prod.yaml 
directory=/data/elastalert
user=root
numprocs=1
priority=1
autostart=true
startsecs=10
satrtretries=3
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=2
stdout_logfile=/tmp/es_alert.log
